# packages/sysutils/systemd/config/system.d/tmp-emulation.mount
[Unit]
Description=Emulation RAM disk
Requires=network-online.service
After=network-online.service
Before=kodi.service
[Mount]
What=tmpfs
Where=/tmp/emulation
Options=size=256m
Type=tmpfs
[Install]
WantedBy=multi-user.target
